# webserver_shutdown.conf
start on started pkgctl-py3k
stop on runlevel [06]
exec /usr/local/bin/python3 /usr/webserver_shutdown.py
